# systemd/flippa.service
[Unit]
Description=flippa Web Application
After=network.target

[Service]
Type=simple
User=flippa
Group=flippa
WorkingDirectory=/opt/flippa
ExecStart=/usr/bin/perl /opt/flippa/app.pl daemon -l http://*:3000
Restart=always
RestartSec=5

# Environment
Environment=FLIPPA_CONFIG=/opt/flippa/config.yml
Environment=MOJO_MODE=production

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/flippa/data /opt/flippa/log
PrivateTmp=true

[Install]
WantedBy=multi-user.target
