# systemd/flippa-prices.service
[Unit]
Description=flippa Price Updater
After=network.target

[Service]
Type=oneshot
User=flippa
Group=flippa
WorkingDirectory=/opt/flippa
ExecStart=/usr/bin/perl /opt/flippa/update_prices.pl --verbose

# Environment
Environment=FLIPPA_CONFIG=/opt/flippa/config.yml

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/flippa/data /opt/flippa/log
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=flippa-prices
